# deis-cleanup.timer
[Unit]
Description=deis-cleanup
Requires=docker.service

[Timer]
OnCalendar=daily

[X-Fleet]
Global=true
